# cron.daily/clean-old-bootcharts
#!/bin/bash
if [ -d /var/log/bootchart ]; then
	find /var/log/bootchart -type f -mtime +14 -exec rm {} \;
fi